# installer/config/supervisor/data/etc/supervisord/backoffice.conf
[program:backoffice]
process_name=%(process_num)s
directory=/data/src
environment=PYTHONPATH="/data/bin:/data/lib:/data/src:$PYTHONPATH"
command=/usr/local/bin/python2.7 -m backoffice.main
    --port=%(process_num)s
    --debug=False
    --track=True
    --prefork_process=8
    --log_file_prefix=/data/var/log/supervisord/%(program_name)s-%(process_num)s.log
numprocs=8
numprocs_start=9000
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stopwaitsecs=600
killasgroup=true
priority=999
user=supervisord
