# installer/config/supervisor/data/etc/supervisord/tasks.conf
[program:tasks]
directory=/data/src
environment=PYTHONPATH="/data/bin:/data/lib:/data/src:$PYTHONPATH"
command=/usr/local/bin/celery worker -A tasks.main -c 8 --loglevel=INFO
stdout_logfile=/data/var/log/celeryd/worker-tasks-main.log
stderr_logfile=/data/var/log/celeryd/worker-tasks-main-error.log
numprocs=8
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stopwaitsecs=600
killasgroup=true
priority=998
user=supervisord
